# systemd/referral_bot.service
[Unit]
Description=Telegram Referral Bot
After=network.target

[Service]
User=your_username
WorkingDirectory=/path/to/your/bot
ExecStart=/usr/bin/python3 /path/to/your/bot/referral_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
